# deploy/issa-api.service
[Unit]
Description=ISSA Trainer Sync API (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=issa
WorkingDirectory=/opt/issa-bot/api
# .env с BOT_TOKEN — общий с ботом (лежит в /opt/issa-bot/.env)
EnvironmentFile=/opt/issa-bot/.env
# слушаем только localhost — наружу отдаёт Caddy на пути /api/*
ExecStart=/opt/issa-bot/.venv/bin/uvicorn api:app --host 127.0.0.1 --port 4100
Restart=on-failure
RestartSec=3
Group=issa

# ── sandbox / hardening ──
# Пишем ТОЛЬКО в каталог api (state.db + WAL + backups/).
ProtectSystem=strict
ReadWritePaths=/opt/issa-bot/api
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# API слушает только localhost; наружу — Caddy. Нужны TCP/UNIX-сокеты:
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
